# dockerfiles/tallyuser-authentication-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallyuser-authentication-servo]
command=/tallycheck/tallyuser-authentication-servo/bin/start.sh
directory=/tallycheck/tallyuser-authentication-servo/conf