# dockerfiles/tallyadmin-authority-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallyadmin-authority-servo]
command=/tallycheck/tallyadmin-authority-servo/bin/start.sh