# dockerfiles/tallybus-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallybus-servo]
command=/tallycheck/tallybus-servo/bin/start.sh