# dockerfiles/tallybiz-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallybiz-servo]
command=/tallycheck/tallybiz-servo/bin/start.sh