# tallyuser/tallyuser-servo/dockerfiles/tallyuser-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallyuser-servo]
command=/tallycheck/tallyuser-servo/bin/start.sh