# dockerfiles/tallyadmin-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallyadmin-servo]
command=/tallycheck/tallyadmin-servo/bin/start.sh