# dockerfiles/tallymanagement-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallymanagement-servo]
command=/tallycheck/tallymanagement-servo/bin/start.sh