# dockerfiles/tallybiz-servo-4-admin/supervisord.conf
[supervisord]
nodaemon=true

[program:tallybiz-servo-4-admin]
command=/tallycheck/tallybiz-servo-4-admin/bin/start.sh