# dockerfiles/tallybus-servo-4-admin/supervisord.conf
[supervisord]
nodaemon=true

[program:tallybus-servo-4-admin]
command=/tallycheck/tallybus-servo-4-admin/bin/start.sh