# dockerfiles/tallyuser-servo/supervisord.conf
[supervisord]
nodaemon=true

[program:tallyuser-servo-4-admin]
command=/tallycheck/tallyuser-servo-4-admin/bin/start.sh